# com.james.dashboard-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.james.dashboard-autopush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/james/Claude/market-research-dashboard/sync-and-push.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/james/Claude/market-research-dashboard/auto-push.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/james/Claude/market-research-dashboard/auto-push.log</string>
</dict>
</plist>
